# systemd/cobertura-progress-watchdog.timer
[Unit]
Description=Cobertura MSP - revisar estancamiento cada 5 minutos

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=cobertura-progress-watchdog.service

[Install]
WantedBy=timers.target
